# cron_jobs.cron
# run the script daily

* * * * * /home/ubuntu/bash_security/analyze_suspicious.sh >> /home/ubuntu/bash_security/cron_log.log
