# cron_jobs.cron
# run the script daily

10 18 * * * /home/ubuntu/bash_security/analyze_suspicious.sh >> /home/ubuntu/bash_security/cron_log.log
