# my_cron_jobs.cron
# run the script daily

32 01 * * * /home/ubuntu/bash_security/analyze_suspicious.sh >> /home/ubuntu/bash_security/cron_log.log
